# cli/packages/cli/src/generators/firewall/scripts/writ-refresh.timer
[Unit]
Description=Writ firewall refresh timer (every 5 minutes)
Documentation=https://github.com/writ/writ

[Timer]
# Run shortly after boot, then every 5 minutes thereafter. AccuracySec
# allows systemd to coalesce with other timers to reduce wakeups.
OnBootSec=30s
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=writ-refresh.service

[Install]
WantedBy=timers.target
